# reolink-unifi-bridge/rootfs/etc/services.d/mosquitto/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Service: mosquitto (Internal MQTT Broker)
# ==============================================================================

bashio::log.info "Starting Mosquitto MQTT broker..."

exec /usr/sbin/mosquitto \
    -c /etc/mosquitto/conf.d/local.conf
