# reolink-unifi-bridge/rootfs/etc/services.d/dashboard/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Service: dashboard (Status Web UI on port 8099)
# Starts immediately - shows "starting" state while other services come up.
# ==============================================================================

DASHBOARD_PORT=$(jq -r '.dashboard' /tmp/actual-ports.json)
bashio::log.info "Starting Reolink Bridge Dashboard on port ${DASHBOARD_PORT}..."

# Export actual ports from 05-port-selection.sh for the Node.js app
NEOLINK_PORT=$(jq -r '.neolink'     /tmp/actual-ports.json)
GO2RTC_PORT=$(jq  -r '.go2rtc_rtsp' /tmp/actual-ports.json)
GO2RTC_API=$(jq   -r '.go2rtc_api'  /tmp/actual-ports.json)
ONVIF_BASE=$(jq   -r '.onvif_base'  /tmp/actual-ports.json)
export NEOLINK_PORT GO2RTC_PORT GO2RTC_API ONVIF_BASE DASHBOARD_PORT
export OPTIONS_FILE="/data/options.json"

cd /opt/dashboard || exit 1

exec node server.js
