# reolink-unifi-bridge/rootfs/etc/services.d/dashboard/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Service: dashboard (Status Web UI on port 8099)
# Starts immediately - shows "starting" state while other services come up.
# ==============================================================================

bashio::log.info "Starting Reolink Bridge Dashboard on port 8099..."

# Export config values as environment variables for the Node.js app
NEOLINK_PORT=$(bashio::config 'neolink_port')
export NEOLINK_PORT
GO2RTC_PORT=$(bashio::config 'go2rtc_port')
export GO2RTC_PORT
export DASHBOARD_PORT=8099
export MQTT_HOST="127.0.0.1"
export MQTT_PORT=1883
export OPTIONS_FILE="/data/options.json"

cd /opt/dashboard || exit 1

exec node server.js
