# reolink-unifi-bridge/rootfs/etc/services.d/neolink/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Service: neolink (Reolink Baichuan → RTSP bridge)
# ==============================================================================

# Kill any stale neolink from a previous run (host_network: port stays bound)
pkill -x neolink 2>/dev/null && sleep 1 || true

CONFIG_FILE="/data/neolink/neolink.toml"
LOG_LEVEL=$(bashio::config 'log_level')

# Map HA log level to Neolink log level
case "${LOG_LEVEL}" in
    "error") NEOLINK_LOG="error" ;;
    "warn")  NEOLINK_LOG="warn" ;;
    "info")  NEOLINK_LOG="info" ;;
    "debug") NEOLINK_LOG="debug,neolink=debug,neolink_core=info" ;;
    *)       NEOLINK_LOG="info" ;;
esac

bashio::log.info "Starting Neolink (log level: ${NEOLINK_LOG})..."
bashio::log.info "Config: ${CONFIG_FILE}"

export RUST_LOG="${NEOLINK_LOG}"
# GStreamer error logging (level 2 = warnings+errors only, 3 = info)
export GST_DEBUG_NO_COLOR=1
export GST_DEBUG=2

exec neolink rtsp \
    --config "${CONFIG_FILE}"
