# reolink-unifi-bridge/rootfs/etc/services.d/go2rtc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Service: go2rtc (RTSP Proxy / Restreamer)
# ==============================================================================

# go2rtc connects directly to cameras via reolink:// protocol
# (no dependency on neolink RTSP server anymore)
# neolink runs independently for optional MQTT motion detection

CONFIG_FILE="/data/go2rtc/go2rtc.yaml"

# Kill any stale go2rtc from a previous run (host_network: port stays bound after restart)
pkill -x go2rtc 2>/dev/null && sleep 1 || true

bashio::log.info "Starting go2rtc..."

exec go2rtc \
    -config "${CONFIG_FILE}"
