# reolink-unifi-bridge/rootfs/etc/services.d/go2rtc/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Service: go2rtc (RTSP Proxy / Restreamer)
# ==============================================================================

NEOLINK_PORT=$(bashio::config 'neolink_port')

bashio::log.info "Waiting for Neolink RTSP to be ready on port ${NEOLINK_PORT}..."
wait-for-port 127.0.0.1 "${NEOLINK_PORT}" 60 || {
    bashio::log.warning "Neolink not ready after 60s, starting go2rtc anyway..."
}

CONFIG_FILE="/data/go2rtc/go2rtc.yaml"

bashio::log.info "Starting go2rtc..."

exec go2rtc \
    -config "${CONFIG_FILE}"
